# ssl/lets.cron
0 6 1 2,4,6,8,10,12 * /root/lets-renew.sh &>/root/lets-renew.log
